# relayguard/deploy/relayguard-agent.service
[Unit]
Description=RelayGuard Agent 转发节点
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/relayguard-agent -data /etc/relayguard-agent
Restart=always
RestartSec=3
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
